# code/etc/systemd/system/plymouth-quit.service.d/exit0-override.conf
[Service]
ExecStart=
ExecStart=-/bin/sh -c "/bin/plymouth quit; exit 0"
