# code/etc/update-motd.d/00-header

Ubuntu 16.04.6 LTS (xenial)
University of Minnesota, Computer Science & Engineering
Questions, Comments, Problems?
Email operator@cs.umn.edu or call 612-625-0876


